# deploy/krs-scanner.service
[Unit]
Description=KRS Batch Scanner
After=network.target

[Service]
Type=simple
User=worker
WorkingDirectory=/opt/rdf-api-project
EnvironmentFile=/opt/rdf-api-project/.env
ExecStart=/opt/rdf-api-project/.venv/bin/python -m batch.runner
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
