# deploy/metadata-backfill.service
[Unit]
Description=RDF Metadata Backfill Worker
After=network.target

[Service]
Type=simple
User=worker
WorkingDirectory=/opt/rdf-api-project
EnvironmentFile=/opt/rdf-api-project/.env
ExecStart=/opt/rdf-api-project/.venv/bin/python -m batch.metadata_runner
Restart=on-failure
RestartSec=60

[Install]
WantedBy=multi-user.target
